# user/backup-pwsafe.timer
[Unit]
Description=Run backup-pwsafe.service every day

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
