# user/linconnect.service
[Unit]
Description=LinConnect Server

[Service]
ExecStart=/usr/bin/linconnect-server

[Install]
WantedBy=default.target
