# user/rsync-backup.service
[Unit]
Description=Backup home to backup drive

[Service]
Type=oneshot
ExecStart=/home/%u/bin/rsync-backup
