# user/rsync-backup.service
[Unit]
Description=Backup home to backup drive

[Service]
Type=oneshot
ExecStart=/usr/bin/bash %h/bin/rsync-backup
